# actividad3/service/saludo_service.service
[Unit]
Description=Servicio de saludo SO1

[Service]
Type=simple
ExecStart=/usr/local/bin/saludo_script.sh

[Install]
WantedBy=default.target
